# profile.d/more_paths.sh
export PATH=$PATH:/home/sam/Applications/google_appengine:/home/sam/bin