# profile.d/qmake.sh
#!/bin/bash

export QMAKE=/usr/bin/qmake-qt4
